# deployment/systemd/llm-config-healthcheck.service
[Unit]
Description=LLM Config Manager - Health check monitoring
Documentation=https://github.com/llm-devops/llm-config-manager
After=llm-config-manager.service
Requires=llm-config-manager.service

[Service]
Type=oneshot
User=llmconfig
Group=llmconfig

# Health check command
ExecStart=/usr/bin/curl -f -s -o /dev/null http://localhost:8080/health || /bin/systemctl restart llm-config-manager.service

# Environment
Environment=RUST_LOG=warn

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true

# Resource limits
MemoryMax=64M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=llm-config-healthcheck

[Install]
WantedBy=multi-user.target
